# browserctl-svc.service
[Unit]
Description=browserctl-svc - Browser Control Service
After=network.target

[Service]
Type=simple
User=dave
Group=dave
WorkingDirectory=/home/dave/code/browserctl/svc

Environment=BROWSERCTL_SVC_PORT=9223
Environment=BROWSERCTL_HTTP_PORT=9224
Environment=BROWSERCTL_PROFILE_DIR=/home/dave/.config/browserctl
Environment=BROWSERCTL_EXT_PATH=/home/dave/workspace/skills/chrome-use/extension

ExecStart=/usr/local/bin/browserctl-svc

Restart=on-failure
RestartSec=5

MemoryMax=1G
MemoryHigh=800M

StandardOutput=journal
StandardError=journal
SyslogIdentifier=browserctl-svc

[Install]
WantedBy=multi-user.target